# crontab
0 5 * * * cd /src && /usr/local/bin/python main.py > /dev/null 2>&1
